# deploy/ubuntu/portfolio-site.service.example
[Unit]
Description=Portfolio Site ASP.NET Core application
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/var/www/portfolio-site/current
ExecStart=/usr/bin/dotnet /var/www/portfolio-site/current/PortfolioSite.dll
Restart=always
RestartSec=5
KillSignal=SIGINT
SyslogIdentifier=portfolio-site
User=www-data
Group=www-data
# Generate /etc/portfolio-site/portfolio-site.env with scripts/create-runtime-env.sh
EnvironmentFile=/etc/portfolio-site/portfolio-site.env

[Install]
WantedBy=multi-user.target
